# rc.local
bash /root/proxy-installer/boot_iptables.sh
bash /root/proxy-installer/boot_ifconfig.sh
ulimit -n 65535
/usr/local/etc/3proxy/bin/3proxy /usr/local/etc/3proxy/3proxy.cfg &